# Yocto_files/meta-robot/recipes-robot/robot-server/files/robot-server.service
[Unit]
Description=Robot Vacuum HTTP Server
After=sound.target network.target
Wants=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/robot-server
WorkingDirectory=/opt/robot
Restart=on-failure
RestartSec=5
SupplementaryGroups=audio
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
